# recipes-support/swupdate/swupdate/coral-dev/boot-boot.mount
[Unit]
Description=Boot boot
After=systemd-modules-load.service

[Mount]
What=/dev/mmcblk1p1
Where=/boot/boot
Type=ext4
